# systemd_status_bridge/docker/bridge.service
[Unit]
Description=systemd status bridge
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/app/bin/bridge

[Install]
WantedBy=multi-user.target
